# contrib/karadul-server.service
[Unit]
Description=Karadul Mesh VPN Coordination Server
Documentation=https://github.com/karadul/karadul
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=karadul
Group=karadul
ExecStart=/usr/local/bin/karadul server --config /etc/karadul/server.json --log-level info
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=10s

# Security hardening
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/karadul /etc/karadul
PrivateTmp=yes
PrivateDevices=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=karadul-server

[Install]
WantedBy=multi-user.target
